# autodark_0.1alpha-1/usr/local/bin/autodark/autodark.service
[Unit]
Description=Time-based automatic theme change

[Service]
Restart=on-failure
ExecStart=/usr/local/bin/autodark/.autodark.sh

[Install]
WantedBy=default.target
